# Hazer/fs/lib/udev/rules.d/40-usb_modeswitch-usb730l.rules
# Novatel USB730L
ATTR{idVendor}=="1410", ATTR{idProduct}=="9030", RUN+="usb_modeswitch '/%k'"
